# .docker/crontab
0 * * * * root sh /usr/local/bin/test.sh
#0 * * * * root /usr/local/bin/php /usr/local/bin/write-calendar-json.php
# An empty line is required at the end of this file for a valid cron file.
